# config/systemd/bea-heartbeat.service
[Unit]
Description=Based EA Builds — hourly heartbeat post
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.openclaw/workspace/openclaw-agent-basedea
ExecStart=/usr/bin/node scripts/daily_prompt.js
Environment=HOME=%h
TimeoutStartSec=120
