# config/systemd/bea-discovery.timer
[Unit]
Description=Based EA Builds — daily discovery at 06:30 EAT

[Timer]
OnCalendar=06:30
Persistent=true

[Install]
WantedBy=timers.target
